# systemd/madridlive-watchdog.timer
[Unit]
Description=Run MadridLive production watchdog every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=madridlive-watchdog.service

[Install]
WantedBy=timers.target
